# config/supervisord/graphite.conf
[program:carbon-cache]
command=carbon-cache --config /etc/carbon/carbon.conf --debug start
autostart=true
autorestart=true
